# _agent-linux/includes/ogp_agent.init.gentoo
#!/sbin/runscript

# Distributed under the terms of the GNU General Public License v2

# GF: Config is in $agent_dir/Cfg/Config.pm

agent_dir=OGP_AGENT_DIR
agent_user=OGP_USER

depend() {
	need net
}

start() {
	ebegin "Starting Game Server Agent"
        start-stop-daemon --verbose --chdir $agent_dir --start --background --user $agent_user -e PWD="$agent_dir" --exec screen -d -m -t agent -c screenrc -S agent ./agent_run -pidfile agent_run.pid
        eend $? "Failed to start Game Server Agent"
}

stop() {
	ebegin "Stopping Game Server Agent"
	start-stop-daemon --stop --quiet --pidfile $agent_dir/agent_run.pid
	eend $? "Failed to stop Game Server Agent"
}

